# imageroot/systemd/user/hermes-agent-pod@.service
[Unit]
Description=Hermes agent pod %i
PartOf=hermes-agent@%i.target
Before=hermes-agent-openviking@%i.service hermes-agent-hermes@%i.service hermes-agent-gateway@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PODMAN_SYSTEMD_UNIT=%n
WorkingDirectory=%S/state
ExecStartPre=/usr/local/bin/runagent discover-smarthost
ExecStartPre=/usr/local/bin/runagent sync-agent-runtime --agent-id %i
ExecStartPre=-/usr/bin/podman pod rm --force hermes-agent-%i
ExecStart=/usr/bin/podman pod create --name hermes-agent-%i
ExecStop=/usr/bin/podman pod stop --ignore --time 10 hermes-agent-%i
ExecStopPost=-/usr/bin/podman pod rm --force hermes-agent-%i
TimeoutStopSec=70

[Install]
WantedBy=default.target