# imageroot/systemd/user/hermes-socket@.service
[Unit]
Description=Hermes dashboard unix socket sidecar %i
Requires=hermes-pod@%i.service hermes@%i.service
After=hermes-pod@%i.service hermes@%i.service
PartOf=hermes@%i.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-%S/state/environment
EnvironmentFile=-%S/state/agent_%i.env
WorkingDirectory=%S/state
Restart=on-failure
SuccessExitStatus=143 SIGTERM
KillMode=mixed
Type=notify
NotifyAccess=all
ExecStartPre=/usr/bin/install -d -m 0770 %S/state/dashboard-sockets
ExecStartPre=-/usr/bin/rm -f %S/state/dashboard-sockets/agent-%i.sock
ExecStartPre=-/usr/bin/podman rm --force hermes-socket-%i
ExecStart=runagent /usr/bin/podman run \
		--name hermes-socket-%i \
		--replace \
		--rm \
		--sdnotify=conmon \
		--cgroups=no-conmon \
		--pod hermes-pod-%i \
		--volume %S/state/dashboard-sockets:/sockets:z \
		--tz=${TIMEZONE} \
		${HERMES_AGENT_SOCKET_IMAGE} \
		-d -d \
		UNIX-LISTEN:/sockets/agent-%i.sock,fork,unlink-early,mode=0660 \
		TCP-CONNECT:127.0.0.1:9120
ExecStop=/usr/bin/podman stop --ignore --time 10 hermes-socket-%i
ExecStopPost=-/usr/bin/podman rm --force hermes-socket-%i
TimeoutStopSec=70

[Install]
WantedBy=default.target